# infra/systemd/market-sentinel.service
# The always-on alert engine: holds the Alpaca websocket, evaluates rules, pushes
# alerts and hosts the Discord remote control. Depends on the config unit having
# populated the EnvironmentFile. On
# startup it auto-approves the latest candidates.csv the pre-screen wrote.
[Unit]
Description=market-sentinel alert engine (live)
After=network-online.target market-sentinel-config.service
Wants=network-online.target
Requires=market-sentinel-config.service
# Tolerate the odd restart, but treat 5 within 5 min as a crash loop: systemd
# gives up restarting and enters "failed", which fires the OnFailure hook below.
StartLimitIntervalSec=300
StartLimitBurst=5
OnFailure=market-sentinel-alert@%n.service

[Service]
Type=simple
User=msuser
WorkingDirectory=/opt/market-sentinel/app
EnvironmentFile=/etc/market-sentinel/engine.env
ExecStart=/opt/market-sentinel/venv/bin/python -m alertengine --live --headless
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
